# drupaltestbot/debian/etc/cron.d
*/1 * * * *     www-data   [ -x /usr/bin/drush ] && /usr/bin/drush -r /var/lib/drupaltestbot core-cron
